# systemd/rclone-mounts-helper.service.in
[Unit]
Description=rclone-mounts privileged D-Bus helper
Documentation=man:rclone-mounts-helper(8)
Requires=dbus.socket
After=dbus.socket

[Service]
Type=dbus
BusName=dev.jthoward.RcloneMounts.Helper
ExecStart=@RCLONE_MOUNTS_HELPER_PATH@

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ConfigurationDirectory=rclone-mounts
ReadWritePaths=/etc/systemd/system /etc/credstore.encrypted
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
# systemd-creds seals system-scope credentials to the TPM when present, but
# PrivateDevices= builds a private /dev containing only pseudo-devices — the
# real /dev/tpmrm0 node is never in it, and DeviceAllow= alone only grants
# cgroup-level access, it doesn't make the node visible. BindPaths= mounts
# the actual device node back in; the "-" prefix keeps this a no-op (not a
# unit failure) on hosts with no TPM.
DeviceAllow=/dev/tpmrm0 rw
BindPaths=-/dev/tpmrm0
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
# D-Bus activation only; no WantedBy.
